# build_info/inetutils-telnetd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.telnetd.plist
