# build_info/inetutils-syslogd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.syslogd.plist
