# build_info/inetutils-talkd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.talkd.plist
