# build_info/inetutils-inetd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.inetd.plist
