# build_info/inetutils-ftpd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.ftpd.plist
