# build_info/inetutils-rshd.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/org.gnu.inetutils.rshd.plist
