# build_info/network-cmds.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/com.apple.kdumpd.plist
