# build_info/remote-cmds.conffiles
@MEMO_PREFIX@/Library/LaunchDaemons/com.apple.telnetd.plist
@MEMO_PREFIX@/Library/LaunchDaemons/com.apple.ntalkd.plist
@MEMO_PREFIX@/Library/LaunchDaemons/com.apple.tftpd.plist
