# build_misc/inetutils/org.gnu.inetutils.talkd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

<dict>
    <key>Label</key>
    <string>org.gnu.inetutils.talkd</string>

    <key>POSIXSpawnType</key>
    <string>Interactive</string>
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>Program</key>
    <string>@MEMO_PREFIX@@MEMO_SUB_PREFIX@/libexec/inetutils-talkd</string>

    <key>ProgramArguments</key>
    <array>
        <string>@MEMO_PREFIX@@MEMO_SUB_PREFIX@/libexec/inetutils-talkd</string>
    </array>

    <key>SessionCreate</key>
    <false/>

    <key>Sockets</key>
    <dict>
        <key>TalkdListenerV4</key>
        <dict>
            <key>SockFamily</key>
            <string>IPv4</string>
            <key>SockType</key>
            <string>dgram</string>
            <key>SockProtocol</key>
            <string>UDP</string>
            <key>SockServiceName</key>
            <string>518</string>
        </dict>
    </dict>

    <key>StandardErrorPath</key>
    <string>/dev/null</string>

    <key>inetdCompatibility</key>
    <dict>
        <key>Wait</key>
        <false/>
    </dict>
    <key>ExecuteAllowed</key>
    <true/>
</dict>

</plist>
