# build_misc/profile.d/multiarch.sh
export DYLD_FALLBACK_LIBRARY_PATH=@MEMO_ALL_LIBDIR@:@MEMO_LIBDIR@:${DYLD_FALLBACK_LIBRARY_PATH}:/usr/local/lib:/usr/lib
