# build_misc/profile.d/profile
export PATH='@PATH@'
export PS1='\h:\w \u\$ '
umask 022

export EDITOR=@MEMO_PREFIX@@MEMO_SUB_PREFIX@/bin/editor
export PAGER=@MEMO_PREFIX@@MEMO_SUB_PREFIX@/bin/pager

for i in @MEMO_PREFIX@/etc/profile.d/*.sh ; do
    if [ -r "$i" ]; then
        . $i
    fi
done
